# addons/init/runit/kula/run
#!/bin/sh
# Kula - runit service

exec chpst -u kula:kula /usr/bin/kula --config /etc/kula/config.yaml serve 2>&1
